# scripts/com.uppilot.stack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.uppilot.stack</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/dev-luv/Sites/uppilot/scripts/uppilot-stack-at-login.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/tmp/uppilot-stack.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/uppilot-stack.launchd.err</string>
</dict>
</plist>
